# ttt.conf
[program:ttt]
command=python3 /home/curator2/ttt/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ttt/test.err.log
stdout_logfile=/home/curator2/ttt/test.out.log